# .config/systemd/user/restic-backup.service
[Unit]
Description=Restic backup service

[Service]
ExecStart=~/.config/backup

[Install]
WantedBy=timers.target
