# .config/systemd/user/restic-backup.service
[Unit]
Description=Restic backup service

[Service]
ExecStart=/home/ben/.config/restic/backup

[Install]
WantedBy=timers.target
