# .config/systemd/user/restic-backup.timer
[Unit]
Description=Timer for restic-backup

[Timer]
ONBootSec=2m
OnUnitInactiveSec=1d
Persistent=true

[Install]
WantedBy=basic.target
